# test-swupdate.service
[Unit]
Description=test-swupdate
After=swupdate

[Service]
ExecStart=/run/media/persistent/test-swu.sh

[Install]
WantedBy=multi-user.target
